# ocr.ini
[program:ocr]
directory=/www/ocr
command=/opt/python/bin/gunicorn -c gunicorn.py main:app
process_name=ocr
autostart=false
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/www/ocr/logs/access.log
stderr_logfile=/www/ocr/logs/error.log
